# contrib/systemd/hyprdictate.service.in
# This is a CMake-substituted template. `cmake --install` writes a
# concrete hyprdictate.service into $CMAKE_INSTALL_FULL_DATADIR/../
# systemd/user/ with @HYPRDICTATE_DAEMON_PATH@ resolved to the
# absolute installed hyprdictated path. Non-Nix installers pick that
# unit up with `systemctl --user daemon-reload`.
#
# NixOS consumers do not use this unit directly. The nix-config
# module synthesises its own systemd.user.services.hyprdictate entry
# with ExecStart pointing at the wrapped store path, so PATH-based
# lookup for `wtype` (added by makeWrapper) stays effective.
[Unit]
Description=hyprdictate voice dictation daemon
Documentation=https://github.com/SubiqT/hyprdictate
# Bind to the graphical session so the daemon starts alongside
# Hyprland and PipeWire, and stops cleanly on logout. Users running
# hyprdictate outside a graphical session (e.g. a bare wayland
# compositor session for testing) can override with `systemctl
# --user edit hyprdictate.service`.
PartOf=graphical-session.target
After=graphical-session.target pipewire.service

[Service]
Type=simple
# @HYPRDICTATE_DAEMON_PATH@ is substituted by CMake to the absolute
# path of the installed daemon. %h expands to $HOME.
ExecStart=@HYPRDICTATE_DAEMON_PATH@ --config %h/.config/hyprdictate/config.toml

# Restart on unexpected exit but leave a small backoff so a bad
# config that fails config load doesn't hammer the journal.
Restart=on-failure
RestartSec=2s

# Kill the daemon and its subprocesses (whisper worker threads,
# wtype subprocesses) as a group on stop so a wtype hung on a
# missing Wayland display doesn't outlive the service.
KillMode=mixed
TimeoutStopSec=5s

# Environment variables the daemon or its subprocesses read:
# - WAYLAND_DISPLAY and XDG_RUNTIME_DIR are inherited from the user
#   session (systemd --user propagates them by default).
# - HYPRDICTATE_LOG_LEVEL overrides the daemon's default INFO level.
# Uncomment to crank verbosity when triaging a failed capture:
# Environment=HYPRDICTATE_LOG_LEVEL=debug

[Install]
WantedBy=graphical-session.target
